# services/easyeffects.service
[Unit]
Description=EasyEffects
After=graphical-session.target pipewire.service

[Service]
ExecStart=/usr/bin/easyeffects --gapplication-service
Restart=on-failure

[Install]
WantedBy=default.target
